# autostart/amgbobbycar.desktop
[Desktop Entry]
Name=AMGBobbyCar
Exec=/usr/bin/python3 /home/pi/AMGBobbyCar/amg_bobby_car.py
